# deploy/fikrlovchi-panel.service
[Unit]
Description=Fikrlovchi admin panel (uzluksiz ishlaydigan web server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/fikrlovchi-panel
ExecStart=/usr/bin/node /root/fikrlovchi-panel/src/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
